# debian/dde-appearance/usr/lib/systemd/user/org.deepin.daemon.appearance.service
[Unit]
Description=Deepin Appearance service

[Service]
Type=DBus
BusName=org.deepin.daemon.Appearance
ExecStart=/usr/bin/dde-Appearance
