# misc/systemd/dde-appearance.service
[Unit]
Description=Deepin Appearance service

[Service]
Type=dbus
BusName=org.deepin.dde.Appearance1
ExecStart=/usr/bin/dde-appearance
