# misc/systemd/dde-appearance.service
[Unit]
Description=Deepin Appearance service

[Service]
Type=DBus
BusName=org.deepin.daemon.Appearance1
ExecStart=/usr/bin/dde-appearance
