# misc/systemd/dde-fakewm.service.in
[Unit]
Description=dde window manager

Requires=plasma-kglobalaccel.service
After=plasma-kglobalaccel.service

[Service]
Type=dbus
BusName=com.deepin.wm
ExecCondition=/bin/sh -c 'test "$XDG_SESSION_TYPE" != "wayland" || exit 2'
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-fakewm
Slice=session.slice
